# deployment/systemd/pethub.service
# ──────────────────────────────────────────────────────────────────────────────
# Pethub — Process manager for Bun server
#
# Place at: /etc/systemd/system/pethub.service
# (deployment/scripts/setup-vps.sh does this automatically)
#
# Commands:
#   sudo systemctl start pethub       — start now
#   sudo systemctl enable pethub      — auto-start on boot
#   sudo systemctl restart pethub     — restart (safe, graceful)
#   sudo systemctl status pethub      — check status
#   sudo journalctl -u pethub -f      — tail logs
# ──────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Pethub — PetCare Pro Bun server
Documentation=https://github.com/your-org/pethub
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/pethub

# ── Executable ────────────────────────────────────────────────────────────────
ExecStart=/usr/local/bin/bun run server.ts
ExecReload=/bin/kill -HUP $MAINPID

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=always
RestartSec=5

# ── Environment ────────────────────────────────────────────────────────────────
# Load from a `.env` file that's NOT in the repo (created by GitHub Actions deploy)
EnvironmentFile=/var/www/pethub/.env.production

# ── Security hardening ─────────────────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/www/pethub /var/log/pethub

# ── Resource limits ────────────────────────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
