# deployment/systemd/vethub.service
# ──────────────────────────────────────────────────────────────────────────────
# Vethub — Systemd service unit
#
# Place at: /etc/systemd/system/vethub.service
# (deployment/scripts/setup-vps.sh does this automatically)
#
# Commands:
#   sudo systemctl start vethub      — start now
#   sudo systemctl enable vethub     — auto-start on boot
#   sudo systemctl restart vethub    — restart (safe, graceful)
#   sudo systemctl status vethub     — check status
#   sudo journalctl -u vethub -f     — tail logs
# ──────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Vethub — VetCare Pro Bun server
Documentation=https://github.com/your-org/vethub
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/vethub

# ── Executable ────────────────────────────────────────────────────────────────
ExecStart=/usr/local/bin/bun run server.ts
ExecReload=/bin/kill -HUP $MAINPID

# ── Restart policy ────────────────────────────────────────────────────────────
Restart=always
RestartSec=5

# ── Environment ────────────────────────────────────────────────────────────────
# Load from a `.env` file that's NOT in the repo (created by setup-vps.sh)
EnvironmentFile=/var/www/vethub/.env.production

# ── Security hardening ─────────────────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/www/vethub /var/log/vethub

# ── Resource limits ────────────────────────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
